# baserock_webserver/etc/systemd/system/generate-docs.baserock.org.timer
[Unit]
Description = Generate docs.baserock.org content

[Timer]
OnUnitActiveSec = 10min

[Install]
WantedBy = default.target
