# baserock_webserver/etc/systemd/system/generate-docs.baserock.org.service
[Unit]
Description = Generate docs.baserock.org content

[Service]
Type = oneshot
ExecStart = /home/fedora/generate-docs.baserock.org
User = fedora
Group = fedora

[Install]
WantedBy=default.target
